# conf/server/biostar.supervisor.conf
#
# Settings for the supervisor processes needed by Biostar.
#
[program:biostar]
command=/home/www/sites/biostar-central/live/gunicorn.start.sh
directory=/home/www/sites/biostar-central
user=www
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/www/sites/biostar-central/live/logs/gunicorn.log

[program:worker]
command=/home/www/sites/biostar-central/live/celery.worker.sh
directory=/home/www/sites/biostar-central
user=www
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/www/sites/biostar-central/live/logs/celery-worker.log

[program:beat]
command=/home/www/sites/biostar-central/live/celery.beat.sh
directory=/home/www/sites/biostar-central
user=www
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/www/sites/biostar-central/live/logs/celery-beat.log