# conf/site/test/supervisor.conf

[program:set-forum-env]
user=www
command =/bin/sh /home/www/apps/biostar-engine/conf/site/test/forum_env.sh

[program:set-recipes-env]
user=www
command =/bin/sh /home/www/apps/biostar-engine/conf/site/test/recipes_env.sh


[group:forum]
programs=set-forum-env
user=www
environment=PATH="/home/www/anaconda2/envs/engine/bin:%(ENV_PATH)s",HOME="/home/www"
command=/home/www/anaconda2/envs/engine/bin/uwsgi --ini /home/www/apps/biostar-engine/forum_uwsgi.ini
stderr_logfile=/home/www/apps/biostar-engine/export/logs/supervisor_stderr.log
stdout_logfile=/home/www/apps/biostar-engine/export/logs/supervisor_stdout.log
autostart=true
autorestart=true
stopsignal=QUIT

[group:recipes]
programs=set-recipes-env
user=www
environment=PATH="/home/www/anaconda2/envs/engine/bin:%(ENV_PATH)s",HOME="/home/www"
command=/home/www/anaconda2/envs/engine/bin/uwsgi --ini /home/www/apps/biostar-engine/recipes_uwsgi.ini
stderr_logfile=/home/www/apps/biostar-engine/export/logs/supervisor_stderr.log
stdout_logfile=/home/www/apps/biostar-engine/export/logs/supervisor_stdout.log
autostart=true
autorestart=true
stopsignal=QUIT
