# conf/site/test/test_supervisor.conf
[program:biostar]
user=www
environment=PATH="/home/www/anaconda2/envs/engine/bin:%(ENV_PATH)s",HOME="/home/www"
command=/home/www/anaconda2/envs/engine/bin/uwsgi --ini /home/www/apps/biostar-engine/%(UWSGI_INI)s
stderr_logfile=/export/www/apps/biostar-engine/export/logs/supervisor_stderr.log
stdout_logfile=/export/www/apps/biostar-engine/export/logs/supervisor_stdout.log
autostart=true
autorestart=true
stopsignal=QUIT
