# conf/supervisor.ini
[program:main]
user=www
environment=PATH="/home/www/miniconda3/envs/uwsgi/bin/:%(ENV_PATH)s"
command=/home/www/miniconda3/envs/uwsgi/bin/uwsgi /export/sites/test_engine/conf/uwsgi-main.ini
stderr_logfile=/export/sites/test_engine/export/logs/supervisor-error.log
stdout_logfile=/export/sites/test_engine/export/logs/supervisor-access.log
autostart=true
autorestart=true
stopsignal=QUIT

[program:test]
user=www
environment=PATH="/home/www/miniconda3/envs/uwsgi/bin/:%(ENV_PATH)s"
command=/home/www/miniconda3/envs/uwsgi/bin/uwsgi /export/sites/main_engine/conf/uwsgi-test.ini
stderr_logfile=/export/sites/main_engine/export/logs/supervisor-error.log
stdout_logfile=/export/sites/main_engine/export/logs/supervisor-access.log
autostart=true
autorestart=true
stopsignal=QUIT
