# conf/test/test_supervisor.conf
[program:test]
user=www
environment=PATH="/home/www/miniconda3/envs/test_env/bin/:%(ENV_PATH)s"
command=/home/www/miniconda3/envs/test_env/bin/uwsgi /export/sites/test_engine/conf/test/test_uwsgi.ini
stderr_logfile=/export/sites/test_engine/export/logs/test_super_stderr.log
stdout_logfile=/export/sites/test_engine/export/logs/test_super_stdout.log
autostart=true
autorestart=true
stopsignal=QUIT
